# supervisor-nginx.conf
[program:nginx]
command = /usr/sbin/nginx
user = www-data